# mqtt-chat/api.service
[Unit]
Description=MQTT Chat API
After=network.target mysql.service docker.service
Wants=mysql.service
[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/mqtt-chat
ExecStart=/usr/bin/gunicorn -w 4 -b 127.0.0.1:5000 api:app
Restart=always
RestartSec=10
[Install]
WantedBy=multi-user.target
